# files/etc/crontabs/root
0 */1 * * * sh -c "/usr/sbin/ntpdate time1.cloud.tencent.com;hwclock -w"
